# legacy_v1/deploy/systemd/kakoito-webhook.service
[Unit]
Description=Kakoito Bot (webhook)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/kakoito
EnvironmentFile=/opt/kakoito/.env
ExecStart=/opt/kakoito/.venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
